# config/prod/crontab
* * * * * /config/prod/start_laravel_cron.sh > /dev/null 2>&1
